# sglang.service
[Unit]
Description=SGLang Qwen3-14B LLM Service
Documentation=https://github.com/hiyouga/sglang
After=network.target

[Service]
Type=simple
User=tchen19
WorkingDirectory=/home/tchen19/LLM
Environment="PATH=/home/tchen19/LLM/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="PYTORCH_CUDA_ALLOC_CONF=max_split_size_mb:64"
Environment="TORCH_CUDA_ALLOC_CONF=garbage_collection_threshold:0.6,max_split_size_mb:64"

ExecStart=/home/tchen19/LLM/.venv/bin/python -m sglang.launch_server \
  --model-path Qwen/Qwen3-14B \
  --reasoning-parser qwen3 \
  --host 0.0.0.0 \
  --port 8000 \
  --device cuda \
  --tensor-parallel-size 2 \
  --dtype half \
  --cpu-offload-gb 48 \
  --offload-mode cpu \
  --mem-fraction-static 0.65 \
  --checkpoint-engine-wait-weights-before-ready \
  --max-running-requests 1 \
  --max-queued-requests 2 \
  --log-level info

# Restart policy
Restart=always
RestartSec=10
StartLimitInterval=60s
StartLimitBurst=3

# Resource limits (adjust based on your needs)
LimitNOFILE=65536
LimitNPROC=4096

# Prevent OOM killer from killing this process
OOMScoreAdjust=-999

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sglang

[Install]
WantedBy=multi-user.target
